# kafka/supervisor-kafka.conf
[program:kafka]
command=/opt/kafka_2.8.0-0.8.1.1/kafka.sh
autostart=true
autorestart=true