# kafka/supervisor-zk.conf
[program:zk]
command=/opt/kafka_2.8.0-0.8.1.1/zk.sh
autostart=true
autorestart=true